# conf.d/tigetvnc.conf
[program:vnc]
command=vncserver -depth 24 -geometry 1024x768 :1
autorestart=true
