# conf.d/tigetvnc.conf
[program:vnc]
command=vncserver :1
autorestart=true
