# conf.d/tightvnc.conf
[program:tightvncserver]
command=vncserver -depth 32 -pixelformat rgb888 -geometry 1440x900 :1
autorestart=true
