# conf.d/tightvnc.conf
[program:tightvncserver]
command=vncserver :1
autorestart=true
