# conf.d/tightvnc.conf
[program:tightvncserver]
command=vncserver
autorestart=true
