# conf.d/unity.conf
[program:fluxbox]
command=unity
autorestart=true
