# conf.d/websockify.conf
[program:websockify]
command=/app/conf.d/websockify.sh
autorestart=true
